# cloudflare-template.service
[Unit]
Description=cloudflare ip list updater service

[Service]
ExecStart=<script path>

[Install]
WantedBy=multi-user.target